# deploy/warden-summary.timer
[Unit]
Description=Post the warden daily summary at 21:00

[Timer]
OnCalendar=*-*-* 21:00:00
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
